# tools/launchd/com.ae.codex-archaeology.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ae.codex-archaeology</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tito/lWORKSPACEl/research/neural-memory/tools/codex_archaeology_pick.sh</string>
    </array>

    <!-- Once a day at 04:00 local. Independent of Claude's 60-min reconciliation reviewer. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/tito/.neural_memory/logs/codex-archaeology.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tito/.neural_memory/logs/codex-archaeology.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Applications/Codex.app/Contents/Resources</string>
        <key>HOME</key>
        <string>/Users/tito</string>
        <!-- Lane label for codex_subagent.sh — distinguishes daily-cron runs from
             interactive nm-builder/ae-builder runs in bridge agent registry + output dirs. -->
        <key>NM_LANE</key>
        <string>cron</string>
    </dict>
</dict>
</plist>
