# tools/launchd/com.ae.codex-todo-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ae.codex-todo-sweep</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tito/lWORKSPACEl/research/neural-memory/tools/codex_todo_sweep.sh</string>
    </array>

    <!-- Daily 05:00 — offset from codex-archaeology (04:00) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/tito/.neural_memory/logs/codex-todo-sweep.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tito/.neural_memory/logs/codex-todo-sweep.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Applications/Codex.app/Contents/Resources</string>
        <key>HOME</key>
        <string>/Users/tito</string>
        <key>NM_LANE</key>
        <string>cron</string>
    </dict>
</dict>
</plist>
